# Linux_Management/GLib/Enumerations/User_Directory/profile.d/g_user_directory_desktop.sh
export G_USER_DIRECTORY_DESKTOP=$HOME/Desktop