# Linux_Management/GLib/Enumerations/User_Directory/profile.d/g_user_directory_pictures.sh
export G_USER_DIRECTORY_PICTURES=$HOME/Pictures