# Linux_Management/GLib/Enumerations/User_Directory/profile.d/g_user_directory_downloads.sh
export G_USER_DIRECTORY_DOWNLOADS=$HOME/Downloads