# Linux_Managment/XDG/Specifications/Base_Directory_Specification/profile.d/xdg_config_dirs.sh
export $XDG_CONFIG_DIRS="/etc/config:/etc/xdg"