# Linux_Managment/GLib/Enumerations/User_Directory/profile.d/g_user_directory_documents.sh
export G_USER_DIRECTORY_DOCUMENTS=$HOME/Documents